# cron/crontab
@weekly cd ${HOME}/laedra/config; git pull
@weekly cd ${HOME}/.config; git push adama
