# cron/crontab
@weekly cd ${HOME}/laedra/config; git pull
@weekly cd /home/koral/.config && hg push /home/hg-public/dotfiles-server
