# etc/systemd/system/nginx.service.d/dependency.conf
[Unit]
Requires=fcgiwrap.service
Requires=php-fpm.service
